# deploy/systemd/kmia-daily-workflow.timer
[Unit]
Description=Timer for KMIA Daily Workflow
# NO REAL TRADING EXECUTION

[Timer]
# Run at 07:00, 12:30, and 19:30 local server time
OnCalendar=*-*-* 07:00:00
OnCalendar=*-*-* 12:30:00
OnCalendar=*-*-* 19:30:00
Persistent=true
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
