# deploy/systemd/kmia-health-check.timer
[Unit]
Description=Timer for KMIA System Health Check
# NO REAL TRADING EXECUTION

[Timer]
# Run every 6 hours
OnCalendar=03:00:00,09:00:00,15:00:00,21:00:00
Persistent=true
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
